# supervisord.conf
[supervisord]
nodaemon=true

[program:wrapdocker]
command=wrapdocker

[program:kubernetes]
command=/start.sh
